# config/systemd/user/fastmail-tunnel.service
[Unit]
Description=SSH tunnel for FastMail SMTP
After=network.target

[Service]
ExecStart=/usr/bin/ssh -v -N -L 1090:smtp.fastmail.com:587 mail-tunnel
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
