# config/systemd/user/restash.target
[Unit]
Description=Restash Backup, Maintenance and Check

[Install]
WantedBy=timers.target
