# config/systemd/user/restic.target
[Unit]
Description=Restic Backup, Maintenance and Check

[Install]
WantedBy=timers.target
