# config/systemd/user/restash-notify@.service
[Unit]
Description=Restash failure notification (%i)

[Service]
Type=oneshot
ExecStart=/bin/sh -l -c 'exec restash-notify %i'
